# build/controller-docker/supervisord.conf
[program:nexns-wsgi]
command=/app/.venv/bin/gunicorn main.wsgi:application -b unix:/var/run/nexns/wsgi.sock
directory=/app/controller
user=www-data
group=www-data
process=1
autostart=true
autorestart=true
stderr_logfile=/data/data/log/controller-wsgi/gunicorn.err.log
stdout_logfile=/data/data/log/controller-wsgi/gunicorn.out.log

[program:nexns-asgi]
command=/app/.venv/bin/daphne -u /var/run/nexns/asgi.sock main.asgi:application
directory=/app/controller
user=www-data
group=www-data
process=1
autostart=true
autorestart=true
stderr_logfile=/data/data/log/controller-asgi/daphne.err.log
stdout_logfile=/data/data/log/controller-asgi/daphne.out.log
